# services/lte.service
[Unit]
Description=My LTE Script

[Service]
Type=simple
ExecStart=/root/start_LTE.sh

[Install]
WantedBy=multi-user.target